# supervisor.conf/button_service.conf
[program:button_service]
directory = /home/pi/Vpet_demo
command = python demo_MAIN.py
stdout_logfile = /home/pi/Vpet_demo/demo_MAIN.log
