# supervisor.conf/hello.conf
[program:hello]
user = pi
command = python /home/pi/hello.py
stdout_logfile = /home/pi/hello.log


